# ops/cron/backup-weekly.cron
# Weekly full backup for spiru-ops
# Every Sunday at 03:00
0 3 * * 0 cd /home/stefano/Documents/Projects/spiru-ops && /bin/bash scripts/backup_spiru_ops.sh >> storage/logs/backup_weekly.log 2>&1
