# ops/cron/daily.cron
# spiru-ops — cron entry (production)
# Runs once/day at 06:10 UTC (reduced from 4x/day to cut Brave Search API costs ~75%).
# Install with: crontab -e  (then paste the line below, adjusting USER if needed)
#
# To install: crontab -e
10 6 * * * cd /home/stefano/Documents/Projects/spiru-ops && PROFILE=kb_first /bin/bash /home/stefano/Documents/Projects/spiru-ops/pipelines/cron_run_daily.sh >> /home/stefano/Documents/Projects/spiru-ops/storage/logs/cron_crontab.log 2>&1
