# ops/cron/weekly.cron
0 10 * * 1 cd /home/stefano/Documents/Projects/spiru-ops && ls -1 storage/artifacts/*_report.md | tail -n 7 > storage/artifacts/weekly_inputs.txt
